# files/config/svsm/sv/pulseaudio/run
#!/usr/bin/env svsm-runscript

exec /usr/bin/pulseaudio --log-target=syslog
