# files/config/svsm/sv/kupfer/run
#!/usr/bin/env svsm-runscript

exec /usr/bin/kupfer --no-splash
