# files/config/svsm/sv/update-notifier/run
#!/usr/bin/env svsm-runscript

depends_on awesome

exec /usr/bin/update-notifier
