# files/config/svsm/sv/gnome-polkit/run
#!/usr/bin/env svsm-runscript

exec /usr/lib/policykit-1-gnome/polkit-gnome-authentication-agent-1
