# files/config/svsm/sv/print-applet/run
#!/usr/bin/env svsm-runscript

exec /usr/bin/system-config-printer-applet
