# files/config/svsm/sv/gcdemu/run
#!/usr/bin/env svsm-runscript

exec /usr/bin/gcdemu
