# files/config/svsm/sv/awesome/run
#!/usr/bin/env svsm-runscript

export wpsetters=feh
exec awesome
